# src/framework_cli/template/infra/backup/backup.timer.jinja
[Unit]
Description=Daily encrypted backup of the {{ project_name }} stack

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
